# development/content/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uwsgi]
command=/opt/virtualenv/bin/uwsgi --ini /var/www/uwsgi.ini
